# radio-backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
loglevel=info

; Icecast starts first — the API needs it internally on localhost:8000
[program:icecast]
command=icecast2 -n -c /etc/icecast2/icecast.xml
user=icecast2
autostart=true
autorestart=true
startretries=5
stdout_logfile=/app/logs/icecast.log
stderr_logfile=/app/logs/icecast.log
stdout_logfile_maxbytes=5MB
priority=10

; nginx is the public front door on Railway's $PORT
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.log
stdout_logfile_maxbytes=5MB
priority=20

; FastAPI/uvicorn on internal port 8099
[program:api]
command=python -m uvicorn main:app --host 127.0.0.1 --port 8099 --workers 1 --log-level info --access-log --no-server-header
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api.log
stdout_logfile_maxbytes=5MB
priority=30
